# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/opt/signstatus/start.sh
StandardInput=tty
StandardOutput=tty

